# vnde/autostart/vn-news-panel.desktop
[Desktop Entry]
Type=Application
Name=VN News Panel
Exec=vn-news-panel
X-GNOME-Autostart-enabled=true
NoDisplay=true
