# vnde/autostart/vnde-panel.desktop
[Desktop Entry]
Type=Application
Name=VNDE Panel
Exec=vnde-gnome-panel
OnlyShowIn=GNOME;
X-GNOME-Autostart-enabled=true
NoDisplay=true
